# conf.d/tachidesk.conf
[program:tachidesk]
command=java -jar "/home/suwayomi/startup/tachidesk_latest.jar"
autorestart=true

stdout_logfile=%(ENV_TACHIDESK_SDOUT)s
stdout_logfile_maxbytes=0