# conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 5800 localhost:5900
autorestart=true