# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -passwd %(ENV_VNC_PASSWORD)s -shared
autorestart=true