# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever %(ENV_LAUNCH_VNC_CMD)s -shared
autorestart=true