# conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -auth ~/.Xauthority %(ENV_DPI_OPTIONS)s -screen 0 %(ENV_XFB_SCREEN)s -listen tcp -ac
autorestart=true